# scripts/services/flyMS.service
[Unit]
Description=flyMS
After=network.target robotcontrol.service pru_handler.service flyMS_webserver.service generic-board-startup.service

[Service]
User=debian
Restart=always
PIDFile=/var/run/flyMS.pid
ExecStart=/home/debian/bin/flyMS_app

[Install]
WantedBy=multi-user.target
